# board/raspberrypi/stak/root/etc/udev/rules.d/70-persistent-net.rules
SUBSYSTEM=="net", ACTION=="add", DRIVERS=="?*", ATTR{type}=="1", KERNEL=="eth*", NAME="%k", RUN+="/sbin/ifup eth%n"
SUBSYSTEM=="net", ACTION=="remove", DRIVERS=="?*", ATTR{type}=="1", KERNEL=="eth*", NAME="%k", RUN+="/sbin/ifdown eth%n"